# client/.docker/supervisord/haven.ini
[supervisord]
nodaemon=true

[group:vault]
programs=haproxy,nginx

[program:haproxy]
command=/usr/sbin/haproxy -f /etc/haproxy/haproxy.cfg
autostart=true
autorestart=true
priority=20
logfile=/var/log/haproxy.log
user=haproxy

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
logfile=/var/log/nginx.log
user=nginx